# conf/servcron.pull.timer
[Unit]
Description=Servecron pull timer 

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
